# communitysite.service
[Unit]
Description=Community Website
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/communitysite
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node ./node_modules/next/dist/bin/next start -p ${PORT}
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target