# communitysite.service
[Unit]
Description=Community Website
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/communitysite
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target